# elk/filebeat/filebeat.service
[Unit]
Description=FileBeat Log Collector

[Service]
WorkingDirectory=/home/koqizhao/elk/filebeat
ExecStart=/home/koqizhao/elk/filebeat/filebeat -e
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
