# middleware/monitoring/metric/grafana/grafana.service
[Unit]
Description=Grafana Visualization System

[Service]
User=koqizhao
Group=koqizhao
WorkingDirectory=BASE_DIR
StandardOutput=file:LOG_DIR/grafana.log
StandardError=file:LOG_DIR/grafana.log
ExecStart=/bin/sh ./start.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
